# POC/cc-server/linux-version/caption-viewer.service
[Unit]
Description=Closed Caption TCP Web Viewer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=captionviewer
Group=captionviewer

WorkingDirectory=/opt/caption-viewer

ExecStart=/usr/bin/python3 /opt/caption-viewer/caption_viewer.py

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target